# howto/usb-autosuspend-exception/etc/udev/rules.d/99-usb-power.rules
ACTION=="add", SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_device", RUN+="/usr/local/bin/usb-power-config.sh '%p'"
